# speedtest.service
[Unit]
Description=Speedtest
After=network.target

[Service]
User=root
ExecStart=/usr/bin/python3  /root/speedtest/speedtest.py